# infrastructure/ansible/applications/github-runner/files/github-runner-cleanup.service
[Unit]
Description=GitHub Runner Cleanup Service
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/github-runner-cleanup.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-runner-cleanup
User=root
Group=root

# Resource limits
TimeoutStartSec=30m
CPUQuota=50%
MemoryMax=2G

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/var/log /var/lib/docker /var/cache /home/runner /tmp /var/tmp

[Install]
WantedBy=multi-user.target