# infrastructure/ansible/applications/github-runner/files/github-runner-cleanup.timer
[Unit]
Description=GitHub Runner Cleanup Timer
Requires=github-runner-cleanup.service

[Timer]
# Run daily at 3 AM
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
# Also run on boot if the last run was missed
Persistent=true
# Run immediately if system was powered down during scheduled time
AccuracySec=1h

[Install]
WantedBy=timers.target