# infrastructure/ansible/applications/github-runner/files/github-runner-disk-monitor.timer
[Unit]
Description=GitHub Runner Disk Space Monitor Timer
Requires=github-runner-disk-monitor.service

[Timer]
# Run every 15 minutes
OnBootSec=5min
OnUnitActiveSec=15min
# Add some randomization to prevent all runners from checking at once
RandomizedDelaySec=5min
# Run immediately if system was powered down during scheduled time
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target