# infrastructure/ansible/templates/consul-wan-health.service.j2
[Unit]
Description=Consul WAN Federation Health Check
After=consul.service tailscaled.service
Requires=consul.service
Wants=tailscaled.service
PartOf=consul.service

[Service]
Type=simple
User=consul
Group=consul
ExecStart=/usr/local/bin/consul-wan-health.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/consul/health-check

[Install]
WantedBy=consul.service