# infrastructure/ansible/templates/etcd-health-check.service.j2
[Unit]
Description=etcd Health Check Service
After=etcd.service
Requires=etcd.service
PartOf=etcd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/etcd-health-check.sh check
User=root
Group=root
StandardOutput=journal
StandardError=journal

# Service hardening
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/etcd
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target