# infrastructure/ansible/templates/etcd-health-check.timer.j2
[Unit]
Description=etcd Health Check Timer
Requires=etcd-health-check.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target