# infrastructure/ansible/templates/gpu-health-check.timer.j2
[Unit]
Description=Run GPU Health Check every 5 minutes
Requires=gpu-health-check.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target