# infrastructure/ansible/templates/pgpool-health-check.timer.j2
[Unit]
Description=pgpool-II Health Check Timer
Requires=pgpool-health-check.service

[Timer]
OnBootSec=3min
OnUnitActiveSec=3min
RandomizedDelaySec=15s

[Install]
WantedBy=timers.target